# enc_v2.4.9p3/targets/ramips/mt76x8/openwrt-enc/etc/uci-defaults/52_mgw7688_firewall
#!/bin/sh

. /lib/functions.sh

config_load firewall

# Accept all connection via WAN interface
uci set firewall.@zone[1].input='ACCEPT'
uci commit firewall

uci commit

exit 0
