# enc_v2.4.9p3/targets/ramips/mt76x8/openwrt-enc/usr/lib/opkg/info/mgw7688-os.list
/usr/bin/mgwp_mpd_monitor
/etc/uci-defaults/54_fstab
/usr/bin/mgwp_get_mb_ipv6
/usr/bin/mgwp_classd
/usr/bin/mgwp_sleep
/root/internet-keep-alive/restart-router.sh
/etc/alsa/linein.asound.state
/usr/bin/mgwp_stat
/etc/alsa/mic.asound.state
/etc/alsa/fm.asound.state
/usr/bin/mgwp_reboot
/usr/bin/enc_info
/usr/bin/mgwp_tailscale
/etc/init.d/mgwp7688_gpio
/usr/bin/mgwp_check_disconnect
/etc/init.d/nwcfg
/usr/bin/mgwp_cfg
/usr/bin/mgwp_streaming_task
/usr/bin/mgwp_monitor
/etc/uci-defaults/52_mgw7688_firewall
/usr/bin/mgwp_upgrade
/usr/bin/mgwp_app
/etc/uci-defaults/55_version_tha
/etc/uci-defaults/53_mgw7688_crontab
/etc/mpd.cfg
/root/internet-keep-alive/restart-interface.sh
/usr/bin/mgwp_streaming_task_v2
/root/internet-keep-alive/dns-test.sh
/usr/bin/speedtest
/usr/bin/mgwp_info
/etc/uci-defaults/51_mgw7688_config
/etc/dropbear/authorized_keys
/usr/bin/mgwp_upgrade_v2
/root/internet-keep-alive/internet-keep-alive.sh
/etc/init.d/mgw7688
