# enc_v2.4.9p3/targets/ramips/mt76x8/openwrt-enc/usr/lib/opkg/info/mpd-full.list
/usr/bin/mpd
/etc/mpd.conf
/etc/init.d/mpd
