# enc_v2.4.9p3/targets/ramips/mt76x8/openwrt-enc/usr/lib/opkg/info/pulseaudio-daemon.list
/usr/lib/libpulse.so.0.23.0
/usr/lib/pulse-14.2/modules/module-virtual-surround-sink.so
/usr/lib/pulseaudio/libpulsecore-14.2.so
/usr/lib/pulse-14.2/modules/module-match.so
/etc/init.d/pulseaudio
/usr/lib/pulse-14.2/modules/module-tunnel-sink.so
/usr/lib/libpulse.so.0
/usr/lib/pulse-14.2/modules/module-pipe-source.so
/usr/lib/pulse-14.2/modules/module-http-protocol-unix.so
/usr/lib/pulse-14.2/modules/module-tunnel-source-new.so
/usr/lib/libpulse-simple.so.0.1.1
/usr/lib/libpulse-simple.so.0
/usr/lib/pulseaudio/libpulsecommon-14.2.so
/etc/pulse/daemon.conf
/usr/lib/libcli.so
/usr/lib/pulseaudio/libpulsedsp.so
/usr/lib/pulse-14.2/modules/module-simple-protocol-unix.so
/usr/lib/libraop.so
/usr/lib/librtp.so
/etc/pulse/default.pa
/usr/bin/pulseaudio
/usr/lib/libprotocol-http.so
/usr/lib/pulse-14.2/modules/module-rescue-streams.so
/usr/lib/pulse-14.2/modules/module-default-device-restore.so
/usr/lib/pulse-14.2/modules/module-suspend-on-idle.so
/usr/lib/pulse-14.2/modules/module-mmkbd-evdev.so
/usr/lib/pulse-14.2/modules/module-native-protocol-unix.so
/usr/lib/pulse-14.2/modules/module-filter-apply.so
/usr/lib/pulse-14.2/modules/module-allow-passthrough.so
/usr/lib/libpulse.so
/usr/lib/pulse-14.2/modules/module-intended-roles.so
/usr/lib/pulse-14.2/modules/module-role-ducking.so
/usr/lib/pulse-14.2/modules/module-rtp-send.so
/usr/lib/pulse-14.2/modules/module-always-sink.so
/usr/lib/pulse-14.2/modules/module-cli-protocol-unix.so
/usr/lib/pulse-14.2/modules/module-card-restore.so
/usr/lib/pulse-14.2/modules/module-cli-protocol-tcp.so
/usr/lib/pulse-14.2/modules/module-alsa-sink.so
/usr/lib/pulse-14.2/modules/module-virtual-source.so
/usr/lib/pulse-14.2/modules/module-position-event-sounds.so
/usr/lib/pulse-14.2/modules/module-raop-sink.so
/usr/lib/libprotocol-native.so
/usr/lib/pulse-14.2/modules/module-rtp-recv.so
/usr/lib/pulse-14.2/modules/module-null-sink.so
/usr/lib/libprotocol-cli.so
/usr/lib/pulse-14.2/modules/module-augment-properties.so
/usr/lib/pulse-14.2/modules/module-filter-heuristics.so
/usr/lib/pulse-14.2/modules/module-virtual-sink.so
/usr/lib/pulse-14.2/modules/module-oss.so
/usr/lib/pulse-14.2/modules/module-null-source.so
/usr/lib/pulse-14.2/modules/module-simple-protocol-tcp.so
/usr/lib/pulse-14.2/modules/module-sine-source.so
/usr/lib/pulse-14.2/modules/module-sine.so
/usr/lib/pulse-14.2/modules/module-remap-source.so
/usr/lib/libalsa-util.so
/etc/pulse/system.pa
/usr/lib/pulse-14.2/modules/module-role-cork.so
/usr/lib/pulse-14.2/modules/module-volume-restore.so
/usr/lib/pulse-14.2/modules/module-device-restore.so
/usr/lib/pulse-14.2/modules/module-ladspa-sink.so
/usr/lib/pulse-14.2/modules/module-switch-on-port-available.so
/usr/lib/pulse-14.2/modules/module-combine-sink.so
/usr/lib/pulse-14.2/modules/module-pipe-sink.so
/usr/lib/pulse-14.2/modules/module-always-source.so
/usr/lib/pulse-14.2/modules/module-tunnel-sink-new.so
/usr/lib/liboss-util.so
/usr/lib/pulse-14.2/modules/module-http-protocol-tcp.so
/usr/lib/pulse-14.2/modules/module-tunnel-source.so
/usr/lib/pulse-14.2/modules/module-native-protocol-tcp.so
/usr/lib/pulse-14.2/modules/module-echo-cancel.so
/usr/lib/pulse-14.2/modules/module-cli.so
/usr/lib/pulse-14.2/modules/module-alsa-source.so
/usr/lib/pulse-14.2/modules/module-remap-sink.so
/usr/lib/pulse-14.2/modules/module-combine.so
/usr/lib/pulse-14.2/modules/module-stream-restore.so
/etc/pulse/client.conf
/usr/lib/pulse-14.2/modules/module-switch-on-connect.so
/usr/lib/libpulse-simple.so
/usr/lib/pulse-14.2/modules/module-device-manager.so
/usr/lib/pulse-14.2/modules/module-alsa-card.so
/usr/lib/pulse-14.2/modules/module-loopback.so
/usr/lib/pulse-14.2/modules/module-detect.so
/usr/lib/libprotocol-simple.so
/usr/lib/pulse-14.2/modules/module-native-protocol-fd.so
