# enc_v2.4.9p3/targets/ramips/mt76x8/openwrt-enc/usr/lib/opkg/info/tailscaled.list
/etc/config/tailscale
/usr/sbin/tailscaled
/etc/init.d/tailscale
